# tools/logrotate_boot.service
[Unit]
Description=run logrotate on boot

[Service]
ExecStart=/usr/sbin/logrotate -f /etc/logrotate.conf

[Install]
WantedBy=multi-user.target