# tools/setup_hwclock.service
[Unit]
Description=set current time from hwclock if present

[Service]
ExecStart=/sbin/hwclock -s

[Install]
WantedBy=multi-user.target